# roles/mongo_pbm_backup/templates/mongo-pbm-backup.timer.j2
[Unit]
Description={{ timer_description }}

[Timer]
OnCalendar={{ oncalendar }}
Persistent=true
Unit={{ unit_name }}

[Install]
WantedBy=timers.target
